# docker/redis/redis-server.ini
[program:redis-server]
command=redis-server /etc/redis.conf
user=redis
directory=/data/redis
stdout_logfile=stdout.log
stderr_logfile=stderr.log
