# docker/postgresql/postgresql.ini
[program:postgresql]
command=/usr/pgsql-9.3/bin/postgres -D /data/postgresql
user=postgres
directory=/data/postgresql
stdout_logfile=stdout.log
stderr_logfile=stderr.log
