# docker/mongodb/mongod.ini
[program:mongod]
command=mongod -f /etc/mongod.conf
user=mongod
directory=/data/mongodb
stdout_logfile=stdout.log
stderr_logfile=stderr.log
